# audio-scheduler.service
[Unit]
Description=Audio Scheduler Service
After=network.target sound.target

[Service]
Type=simple
User=istvan
WorkingDirectory=/home/istvan/Dokumentumok/Dev/audio-scheduler
Environment="FLASK_ENV=production"
Environment="PYTHONUNBUFFERED=1"
Environment="XDG_RUNTIME_DIR=/run/user/1000"
Environment="PULSE_RUNTIME_PATH=/run/user/1000/pulse"

# Use bash to run the gunicorn command
# This avoids SELinux permission issues with direct venv/bin execution
# IMPORTANT: -w 1 (single worker) is critical to prevent duplicate schedulers
ExecStart=/bin/bash -c 'source venv/bin/activate && exec gunicorn \
    -w 1 \
    -b 0.0.0.0:5000 \
    --timeout 120 \
    --log-level info \
    --access-logfile logs/gunicorn_access.log \
    --error-logfile logs/gunicorn_error.log \
    --capture-output \
    wsgi:app'

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=audio-scheduler

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
